# mkosi-crypt/mkosi.extra/usr/lib/systemd/system/systemd-repart.service.d/systemd-repart-gen_tabs.conf
[Unit]
RequiresMountsFor=/run
ConditionPathIsReadWrite=/run

[Service]
ExecStart=
ExecStart=systemd-repart --dry-run=no --key-file=/run/fscrypt.sock --generate-fstab=/run/fstab --generate-crypttab=/run/crypttab
